# udev/99-adb.rules
ACTION=="bind|remove",ATTR{interface}=="ADB Interface",RUN+="$BASE/scripts/adbforward.sh $USER %E{ACTION} %S %E{DEVPATH}"

